# srcs/nginx/srcs/supervisor.ini
[supervisord]
nodaemon=true
user=root

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
startsecs=0
redirect_stderr=true

[program:php-fpm]
command=/usr/sbin/php-fpm7
autostart=true
autorestart=true

[program:ssh]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
